# sw/qsdk/package/KeyWest/kw-base-files/files/etc/init.d/KWN_defaults
#!/bin/sh /etc/rc.common
START=15
USE_PROCD=1

start_service() {
    procd_open_instance

    eth_timer=$(uci get ethernet.ethernet.ethtimer)
    tx_power_radio1=$(uci get wireless.wifi1.TXPowLim5G)
    tx_power_radio2=$(uci get wireless.wifi0.TXPowLim2G)
    tx_chainmask_radio1=$(uci get wireless.wifi1.txchainmask)
    rx_chainmask_radio1=$(uci get wireless.wifi1.rxchainmask)
    shortgi_radio1=$(uci get wireless.@wifi-iface[1].shortgi)
    distance_radio1=$(uci get wireless.wifi1.distance)
    hidden_radio0=$(uci get wireless.@wifi-iface[0].hidden)
    hidden_radio1=$(uci get wireless.@wifi-iface[1].hidden)
    dgackdisable=$(uci get wireless.@wifi-iface[1].dgackdisable)
    
    if [ "$eth_timer" == "" ]
    then
        uci set ethernet.ethernet.ethtimer='0'
    fi
    if [ "$tx_power_radio1" == "" ]
    then
        uci set wireless.wifi1.TXPowLim5G='26'
    fi
    if [ "$tx_power_radio2" == "" ]
    then
        uci set wireless.wifi0.TXPowLim2G='20'
    fi
    if [ "$tx_chainmask_radio1" == "" ]
    then
        uci set wireless.wifi1.txchainmask='3'
    fi
    if [ "$rx_chainmask_radio1" == "" ]
    then
        uci set wireless.wifi1.rxchainmask='3'
    fi
    if [ "$shortgi_radio1" == "" ]
    then
        uci set wireless.@wifi-iface[1].shortgi='1'
    fi
    if [ "$distance_radio1" == "" ]
    then
        uci set wireless.wifi1.distance='5'
    fi
    if [ "$hidden_radio1" == "" ]
    then
        uci set wireless.@wifi-iface[1].hidden='0'
    fi
    if [ "$dgackdisable" == "" ]
    then
        uci set wireless.@wifi-iface[1].dgackdisable='0'
    fi
#2.4 Radio Default values
    uci set wireless.wifi0.txchainmask='1'
    uci set wireless.wifi0.rxchainmask='1'
    uci set wireless.wifi0.rate='0'
    if [ "$hidden_radio0" == "" ]
    then
        uci set wireless.@wifi-iface[0].hidden='1'
    fi

    procd_close_instance
}
